# node.conf
description "Running node service"

author  "trevor@sprout.co.ke"
start on runlevel [2345]
stop on runlevel [!2345]

respawn
exec /var/www/ongair.im/public_html/exec.sh